# docker/nginx/cron/crontab
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Letsencrypt certificates
0 */3 */1 * * certbot renew --allow-subset-of-names